# Source_code/shell_file/profile.d/auto_start.sh
#!/bin/bash

/home/new/Desktop/opencv/opencv-3.4.0/build/darknet/boot_start.py                                                                       
